# meta-tizen/recipes-graphics/wayland/weston/weston.target

[Unit]
Description=Weston Compositor
Requires=dbus.socket
AllowIsolate=true
